# scraper/com.youfirstgersh.dailyscraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.youfirstgersh.dailyscraper</string>

  <key>ProgramArguments</key>
  <array>
    <string>__SCRAPER_DIR__/.venv/bin/python</string>
    <string>-m</string>
    <string>youfirst_scraper.run_daily</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__SCRAPER_DIR__</string>

  <!-- 13:00 and 17:00 fires retry handles the 9:00 run missed; a fully
       successful earlier run makes them no-ops via the .last_run guard. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>13</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>17</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <string>__SCRAPER_DIR__/logs/daily_run.log</string>
  <key>StandardErrorPath</key>
  <string>__SCRAPER_DIR__/logs/daily_run.err.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
